# etc/systemd/system/chronyd.service.d/override.conf
[Unit]
Wants=dev-rtc_pps.device dev-gps_pps.device
After=dev-rtc_pps.device dev-gps_pps.device
