# lib/systemd/system/display-clock.service
[Unit]
Description=Display clock

[Service]
Type=simple
ExecStart=/usr/bin/display-clock

[Install]
WantedBy=multi-user.target
