# etc/systemd/chronyd.service.d/override.conf
[Unit]
Wants=dev-pps0.device
After=dev-pps0.device

