# etc/udev/rules.d/99-gps-name.rules
KERNEL=="ttyS4", SUBSYSTEM=="tty", SYMLINK+="gps0", RUN+="/bin/stty -F /dev/ttyS4 speed 115200"


